# deploy/aether-api.service.d/logging.conf
# MV-system-001 (docs/delivery/MANUAL-VERIFICATION-GAPS.json) — DECISION:
# keep this file-based sink rather than switching to
# `StandardOutput=journal` / `StandardError=journal`.
#
# Why: journald was empty for every aether-* unit ("No journal files were
# found") and this file's own lines carried no timestamp — TWO separate
# defects, not one. Switching the sink to journald would (a) require a
# restart to take effect, (b) leave `/var/log/aether/api.log` behind stale
# and empty going forward, silently breaking every existing runbook
# procedure that greps/tails it (docs/delivery/DEPLOYMENT-RUNBOOK.md §4,
# §5 "Verify Deployment", "Troubleshooting"), and (c) is a bigger, riskier
# infra change than this MV fix's mandate ("minimal diff, no scope creep").
# The actual blocker — "window-scoped forensics for incidents/sweeps is
# impossible" — is fixed by adding an ISO-8601 UTC timestamp to every line
# instead (apps/api/logging_config.json, wired via `uvicorn --log-config` in
# start-api.sh): `grep` / `awk` against this same file can now scope by
# timestamp exactly like journalctl's `--since`/`--until` would.
[Service]
StandardOutput=append:/var/log/aether/api.log
StandardError=append:/var/log/aether/api.log
