# deploy/aether-api.service
# Aether API server (FastAPI/Uvicorn).
# Install: sudo cp deploy/aether-api.service /etc/systemd/system/aether-api.service
#          sudo cp deploy/aether-api.service.d/logging.conf /etc/systemd/system/aether-api.service.d/logging.conf
#          sudo systemctl daemon-reload && sudo systemctl restart aether-api
# Logs:    /var/log/aether/api.log (see aether-api.service.d/logging.conf)
#
# MV-system-001 (docs/delivery/MANUAL-VERIFICATION-GAPS.json): this is the
# first time this unit is tracked in version control — it and its logging.conf
# drop-in were previously host-only files (confirmed byte-identical to the
# live /etc/systemd/system/aether-api.service on 2026-07-18). No behavior
# change here; see aether-api.service.d/logging.conf for the rationale on
# keeping the file-based (not journald) sink, and start-api.sh /
# apps/api/logging_config.json for the accompanying ISO-8601 timestamp fix.
[Unit]
Description=Aether API Server (FastAPI/Uvicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/github_repos/aether-job-career-agent
ExecStart=/home/ubuntu/github_repos/aether-job-career-agent/start-api.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
