# deploy/aether-web.service.d/logging.conf
# MV-system-001 (docs/delivery/MANUAL-VERIFICATION-GAPS.json) — same
# file-based-sink decision as deploy/aether-api.service.d/logging.conf; see
# that file's comment for the full rationale. The timestamp fix for this
# unit lives in start-web.sh (gawk prefix), since Next.js's `next start` has
# no built-in log-format configuration hook the way uvicorn does.
[Service]
StandardOutput=append:/var/log/aether/web.log
StandardError=append:/var/log/aether/web.log
